# plugins/default/vulnerabilities/apache2_slow_loris/flag-publish.service
[Unit]
Description=Simple flag publishing service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/vagrant/publish-flag.sh
Restart=always
RestartSec=3s

[Install]
WantedBy=multi-user.target